# config/deploy/templates/shared/sidekiq.service.erb
[Unit]
Description=Sidekiq workers
After=syslog.target network.target
<%- fetch(:sidekiq_processes).times do |i| -%>
Wants=sidekiq@<%= i %>.service
<% end %>

[Service]
Type=oneshot
ExecStart=/bin/true
ExecReload=/bin/kill -TSTP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
RemainAfterExit=true

[Install]
WantedBy=default.target
